# plantillas-conf/supervisord.conf
# ==========================================
# CONFIGURACIÓN BASE (¡NUNCA BORRAR!)
# ==========================================
[supervisord]
nodaemon=true
user=root

# ==========================================
# BLOQUE 1: EL PROXY (Nginx)
# -> Borra esto si tu app no usa Nginx y responde directo a internet.
# ==========================================
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ==========================================
# BLOQUE 2: LA APLICACIÓN PRINCIPAL (Backend / API)
# -> Solo necesitas cambiar el "command=" según tu lenguaje.
# ==========================================
[program:backend-app]

# Ejemplos de comandos (descomenta el que uses y borra el resto):
# command=java -jar /app/app.jar
# command=uvicorn main:app --host 127.0.0.1 --port 5000
# command=node server.js
command=java -jar /app/app.jar --spring.datasource.url="jdbc:h2:tcp://localhost:9092/inventario"

autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ==========================================
# BLOQUE 3: BASE DE DATOS EMBEBIDA (Ej: H2)
# -> Borra este bloque si usas una BD externa (Postgres) en otro contenedor.
# ==========================================
[program:database]
command=java -cp /app/h2-server.jar org.h2.tools.Server -tcp -tcpPort 9092 -tcpAllowOthers -ifNotExists -baseDir /app/data
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ==========================================
# BLOQUE 4: TAREAS EN SEGUNDO PLANO (Workers / Scripts)
# -> Borra este bloque si tu proyecto no tiene scripts corriendo de fondo.
# ==========================================
[program:worker]
command=python /app/worker_tareas.py
autostart=true
autorestart=true
startretries=3  # Le da 3 intentos para arrancar si falla al inicio
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
